# src/gravity-fingerprints/dbus-com.ispirata.Hemera.Fingerprints.service.in
[Unit]
Description=Hemera Fingerprints Daemon

[Service]
Type=dbus
BusName=com.ispirata.Hemera.Fingerprints
ExecStart=@FULL_INSTALL_BIN_DIR@/gravity-fingerprints

[Install]
Alias=com.ispirata.Hemera.Fingerprints.service
